# assets/linux/alpinezen.service
[Unit]
Description=Breathe life into your workspace with dynamic wallpapers

[Service]
Type=simple
ExecStart=/usr/bin/alpinezen
StartLimitInterval=5
StartLimitBurst=10
Restart=always
RestartSec=10s

[Install]
WantedBy=graphical-desktop.target
